# packaging/rprocd.system.service
[Unit]
Description=rproc background sampler
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/rproc --daemon
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target